# scripts/binkyrailways.service
[Unit]
Description=BinkyRailways
After=network.target
After=dnsmasq.service
Requires=dnsmasq.service

[Service]
ExecStart=/home/pi/binky-server \
  --storage-path=/home/pi/binky-storage \
  --certificate-path=/home/pi/binky-cert.json \
  --published-host=${PUBLISHED_HOST_IP} \
  --published-hostname=${PUBLISHED_HOST_DNS} \
  --http-port=80
EnvironmentFile=/home/pi/binky-server.env
Restart=on-failure

[Install]
WantedBy=multi-user.target