# scripts/binkyrailways.service
[Unit]
Description=BinkyRailways
After=network.target
After=dnsmasq.service
Requires=dnsmasq.service

[Service]
Type=service
ExecStart=/home/pi/binky-server --prom-config-path=/home/pi/prometheus.yml --storage-path=/home/pi/binky-storage

[Install]
WantedBy=multi-user.target